# config/systemd/kun-health.timer
[Unit]
Description=Run Kun health check every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
